# docker/containers/application/etc/supervisor/production.conf
[supervisord]
logfile = /var/www/project/docker/containers/application/var/log/supervisor/supervisor.log
pidfile = /var/www/project/docker/containers/application/var/log/supervisor/supervisor.pid

[program:queue]
command = php /var/www/project/html/framework/yii queue/listen 5 --verbose=1 --isolate=1 --color=0
autostart = true
autorestart = true
stdout_logfile = /var/www/project/docker/containers/application/var/log/supervisor/queue.out.log
stderr_logfile = /var/www/project/docker/containers/application/var/log/supervisor/queue.err.log
