# docker/containers/application/var/spool/cron/config/crontab
LANG=ru_RU.UTF-8

@reboot /usr/bin/supervisord -c /etc/supervisor/application.conf

# * * * * * /usr/local/bin/php /var/www/html/framework/yii
